# myproject.service
[Unit]
Description=uWSGI instance to serve myproject
After=network.target

[Service]
User=wandx
Group=www-data